# tools/macshield/com.qinnovates.macshield.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.qinnovates.macshield</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/macshield</string>
        <string>--trigger</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Library/Preferences/SystemConfiguration/com.apple.airport.preferences.plist</string>
        <string>/Library/Preferences/SystemConfiguration/preferences.plist</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/tmp/macshield.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/macshield.stderr.log</string>

    <!-- NOTE: install.sh replaces /tmp paths above with ~/Library/Logs/ at install time -->
</dict>
</plist>
